# hot-pot/surge-hotspot.service
[Unit]
Description=Proxy WiFi hotspot (routes client traffic through surge-tunnel)
Documentation=file:///home/zhang/Projects/auto-switch-ssh-proxy/hot-pot/README.md
Wants=surge-tunnel.service xray.service NetworkManager.service
After=surge-tunnel.service xray.service NetworkManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/hotspot-up.sh
ExecStop=/usr/local/bin/hotspot-down.sh
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
